# config/studiobrain/systemd/studio-brain-discord-relay.timer
[Unit]
Description=Periodic fallback Studio Brain Discord relay sync

[Timer]
OnBootSec=2m
OnUnitActiveSec=1m
AccuracySec=15s
Persistent=true
Unit=studio-brain-discord-relay.service

[Install]
WantedBy=timers.target
