# config/studiobrain/systemd/studio-brain-disk-alert.service
[Unit]
Description=Studio Brain disk usage check

[Service]
Type=oneshot
ExecStart=/bin/bash /usr/local/bin/studio-brain-disk-alert.sh
