# config/studiobrain/systemd/studio-brain-reboot-watch.service
[Unit]
Description=Studio Brain reboot-required watch
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/bin/bash /usr/local/bin/studio-brain-reboot-watch.sh
