# config/studiobrain/systemd/user/studio-brain-memory-ops-supervisor.timer
[Unit]
Description=Keep Studio Brain memory ops supervisor active

[Timer]
OnBootSec=90s
OnUnitActiveSec=5m
Unit=studio-brain-memory-ops-supervisor.service

[Install]
WantedBy=timers.target
